# logwise-orchestrator/cron/crontab
# Every 1 minute - POST
*/1 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/api/v1/component/sync || echo "$(date -Is) A failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/update-logs-delete-policy || echo "$(date -Is) B failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/sync-livelogs-host || echo "$(date -Is) C failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/api/v1/metric/object-store/sync-delay || echo "$(date -Is) D failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/monitor-spark-job || echo "$(date -Is) E failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/scale-kafka-partitions || echo "$(date -Is) F failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/clean-spark-compact-file || echo "$(date -Is) G failed" >&2

*/5 * * * * curl -fsS --retry 3 --retry-connrefused --max-time 15 \
  -X POST -H "Content-Type: application/json" -H "X-Tenant-Name: dream11" \
  http://app:8080/api/v1/grafana/log || echo "$(date -Is) H failed" >&2